# crontab
0 3 * * * root /backup.sh
0 4 * * * root /update.sh

