# crontab
*/10 * * * * root /cronjob.sh
0 3 * * * root /update.sh

